# custom_script
#/bin/bash

### BEGIN INIT INFO
# Provides: chillispot et freeradius dans le chroot
# Required-Start: $local_fs $network
# Required-Stop: $local_fs $remote_fs
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Mise en place des scripts de démarrage
# Description: Ensemble des scripts custom à lancer au démarrage de la machine
# or wireless LAN access point controller.
### END INIT INFO

#AUTEUR: S0bek
#DESCRIPTION: Script permettant d'appeler les autres scripts à exécuter au démarrage
#DATE: 27/10/2015

LOGPATH="/var/log"
LOGFILE="$LOGPATH/custom_scripts.log"

log_purge () {
    if [ -f $LOGFILE ]; then
        rm -f $LOGFILE
        if [ $? = 0 ]; then echo "***Suppression des logs termine***"; fi
    fi
}

exec_custom_scripts () {

    #Script 1: Démarrage de hostapd
    /etc/init.d/scripts/APSetup.sh | tee -a $LOGFILE

    #Script 2: Démarrage du daemon DHCP juste après la configuration de hostapd
    /etc/init.d/scripts/DHCPSetup.sh | tee -a $LOGFILE

    #Script 3: Vérification du bon foncitonnement du point d'accès
    /etc/init.d/scripts/APStatus.sh | tee -a $LOGFILE

    #Script 4: Dans le sous-dossier scripts/ toujours

}

log_purge && exec_custom_scripts

